# tools/service/web-redshore.service
[Unit]
Description=RedShoreTech WebSite

[Service]
User=root
WorkingDirectory=/root/sites/redshore/
ExecStart=/bin/bash -c 'hugo server --bind=wumaseman.com -p 80 -D -F'
# optional items below
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target

